# packages/sys-apps/rng-tools/files/systemd/rngd.service
[Unit]
Description=rng-tools rngd

[Service]
Restart=always
CapabilityBoundingSet=CAP_SYS_ADMIN
EnvironmentFile=/etc/conf.d/rngd.conf
ExecStart=/usr/sbin/rngd -f -n1 ${STEP} ${WATERMARK} ${DEVICE}

[Install]
WantedBy=basic.target
